# deploy/smc-trader.service
[Unit]
Description=SMC Paper Trading Engine
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/smcbot/smc-trader
EnvironmentFile=/home/smcbot/smc-trader/.env
ExecStart=/home/smcbot/smc-trader/.venv/bin/gunicorn app:app --bind 0.0.0.0:5001 --workers 1 --threads 4 --timeout 120
Restart=always
RestartSec=10
KillSignal=SIGTERM
TimeoutStopSec=30

# Resource limits (Oracle Cloud Free Tier: 4 OCPU, 24GB RAM)
# Leave room for OS and 15+ strategies in future
LimitNOFILE=65535

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=smc-trader

[Install]
WantedBy=multi-user.target
